# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# ── MariaDB ──────────────────────────────────────────────
[program:mysql]
command=mysqld_safe --user=mysql
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/mysql.out.log
stderr_logfile=/var/log/supervisor/mysql.err.log
priority=10

# ── Apache (arranca despues de que MySQL este listo) ──────
[program:apache2]
command=/usr/local/bin/apache-wait-start.sh
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/supervisor/apache2.out.log
stderr_logfile=/var/log/supervisor/apache2.err.log
priority=20
environment=MYSQL_HOST="127.0.0.1",MYSQL_DATABASE="sistema_agencia",MYSQL_USER="agencia_user",MYSQL_PASSWORD="agencia_pass2026"
